# debian/postinst
#!/bin/sh

#set -e

case "$1" in
    configure)

        if [ -f /usr/share/debconf/confmodule ]; then

            . /usr/share/debconf/confmodule
            db_version 2.0

            db_get helixbilling/hostname
            # hostname value now in RET
            . /opt/helixutils/substitute.sh
            echo "Setting $RET as sever_name into nginx config"
            subst /etc/nginx/sites-available/helixbilling "server_name.*" "server_name  $RET;"

            if [ -f /usr/share/dbconfig-common/dpkg/postinst.pgsql  ]; then

                . /usr/share/dbconfig-common/dpkg/postinst.pgsql

                dbc_generate_include='template:/usr/share/python-support/helixbilling/helixbilling/conf/settings.py'
                dbc_generate_include_owner='helixbilling:root'
                dbc_generate_include_perms='640'
                dbc_generate_include_args="-U -o template_infile=/usr/share/helixbilling/settings.py.tpl"

                # you can set the default database encoding to something else
                dbc_pgsql_createdb_encoding="UTF8"

                dbc_go helixbilling $@

            fi
        fi

        if which update-python-modules >/dev/null 2>&1; then
            update-python-modules  helixbilling
        fi

        chown -R helixbilling:helixbilling /var/log/helixbilling
        chmod 750 /var/log/helixbilling

        echo "Installing helixbilling patches"
        sudo -u helixbilling /usr/share/helixbilling/helixbilling_install update

        echo "Installing nginx config"
        chmod 400 /etc/nginx/sites-available/helixbilling-ssl/*.key
        ln -f -s /etc/nginx/sites-available/helixbilling /etc/nginx/sites-enabled
        invoke-rc.d nginx reload

        echo "Installing helixbilling configuration into supervisor"
        ln -f -s /etc/supervisor/process-available/helixbilling.ini /etc/supervisor/process-enabled
        invoke-rc.d supervisor reload

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


